# roles/postgresql_uninstall/files/postgresql.service
[Service]
Environment=PGDATA=/var/lib/pgsql/data/
